# pkg/target/classes/deb/root/etc/systemd/system/jubel.service
[Unit]
Description=Jubel | Packaging
After=syslog.target

[Service]
User=jubel
ExecStart=/opt/jubel/jubel.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
